# conf.d/Mars-websockify.conf
[program:websockify]
command=/app/conf.d/Mars-websockify.sh
autorestart=true
